# internal/chalkctl/templates/chalk-coturn.container.tmpl
# chalk -- coturn media relay (generated by chalkctl init). Host networking;
# plain turn: only (media is DTLS-SRTP E2E regardless).
[Unit]
Description=chalk coturn
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=coturn
Image=docker.io/coturn/coturn:4
Network=host
EnvironmentFile=/etc/chalk/chalk.env
Exec=--log-file=stdout \
  --realm=chalk \
  --use-auth-secret \
  --static-auth-secret=${CHALK_TURN_SECRET} \
  --listening-port=3478 \
  --min-port=49160 \
  --max-port=49200 \
  --fingerprint --no-cli --no-tls --no-dtls

[Service]
Restart=always

[Install]
WantedBy=default.target
