# internal/chalkctl/templates/chalk-update.service.tmpl
# chalk -- weekly auto-update oneshot (generated by chalkctl init).
[Unit]
Description=chalk auto-update run
After=network-online.target

[Service]
Type=oneshot
ExecStart={{ .ChalkctlPath }} update --yes --channel stable
