# internal/chalkctl/templates/chalk-coturn.container.tmpl
# chalk -- coturn media relay (generated by chalkctl init). Host networking;
# alpine image. Settings come from a mounted config file (literal secret, no
# env expansion); LOGGING flags go on the CLI (-n --log-file=stdout) because
# coturn's file logger otherwise swallows output to /var/tmp/turn_*.log and
# `podman logs coturn` stays empty (coturn issues #507/#1528). -n disables the
# config-file log redirection so stdout logging actually reaches the container
# stdout. DETECT_EXTERNAL_IP asks the image's detect-external-ip helper to add
# --external-ip=<detected IPv4> at runtime (IPv4 only; no IPv6).
[Unit]
Description=chalk coturn
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=coturn
Image=docker.io/coturn/coturn:{{ .CoturnTag }}
Network=host
Environment=DETECT_EXTERNAL_IP=yes
Environment=DETECT_RELAY_IP=yes
Volume=/etc/chalk/coturn/turnserver.conf:/etc/coturn/turnserver.conf:ro,Z
Exec=-n --log-file=stdout --simple-log --new-log-timestamp {{ if .TurnVerbose }}--verbose {{ end }}-c /etc/coturn/turnserver.conf

[Service]
Restart=always

[Install]
WantedBy=default.target
