# internal/chalkctl/templates/chalk-update.timer.tmpl
# chalk -- weekly auto-update timer (generated by chalkctl init).
[Unit]
Description=chalk weekly auto-update

[Timer]
# Sunday 04:00 -- after the CI Sunday 02:00 security rebuild (ops-2).
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
