# internal/chalkctl/templates/chalk-postgres.container.tmpl
# chalk -- Postgres {{ .PostgresTag }} (generated by chalkctl init; PG major
# upgrades are MANUAL, never auto).
[Unit]
Description=chalk Postgres
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=postgres
Image=docker.io/library/postgres:{{ .PostgresTag }}
Network=chalk.network
Volume=chalk-pgdata.volume:/var/lib/postgresql/data
EnvironmentFile=/etc/chalk/chalk.env
Environment=POSTGRES_DB=chalk
Environment=POSTGRES_USER=chalk
HealthCmd=pg_isready -U chalk -d chalk
HealthInterval=5s
HealthTimeout=3s
HealthRetries=10

[Service]
Restart=always

[Install]
WantedBy=default.target
