# internal/chalkctl/templates/chalk-coturn.container.tmpl
# chalk -- coturn media relay (generated by chalkctl init). Host networking;
# reads a generated config file (literal secret) via -c, avoiding the env-var
# expansion that a CLI a CLI --static-auth-secret dollar-var would collapse to empty.
[Unit]
Description=chalk coturn
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=coturn
Image=docker.io/coturn/coturn:4
Network=host
Volume=/etc/chalk/coturn/turnserver.conf:/etc/coturn/turnserver.conf:ro,Z
Exec=-c /etc/coturn/turnserver.conf

[Service]
Restart=always

[Install]
WantedBy=default.target
