# internal/chalkctl/templates/chalkd.container.tmpl
# chalk -- chalkd app (generated by chalkctl init). Digest-pinned to the exact
# signed image; the digest is immutable. To update, use `chalkctl update`.
#
#   image  : {{ .Image }}:{{ .Version }}
#   digest : {{ .Digest }}
[Unit]
Description=chalk chalkd
Wants=network-online.target
After=network-online.target chalk-postgres.service
Requires=chalk-postgres.service

[Container]
ContainerName=chalkd
Image={{ .Image }}@{{ .Digest }}
PullPolicy=missing
Network=chalk.network
Volume=chalk-blobs.volume:/var/lib/chalk/blobs
EnvironmentFile=/etc/chalk/chalk.env
Environment=CHALK_LISTEN=:8443
Environment=CHALK_TLS_MODE=off
Environment=CHALK_LOG_FORMAT=json
Environment=CHALK_INSTANCE_ID=instance-1
Environment=CHALK_DB_URL=postgres://chalk:${CHALK_PG_PASSWORD}@postgres:5432/chalk?sslmode=disable
Environment=CHALK_VOICE_ENABLED={{ .VoiceEnabled }}
{{- if .VoiceEnabled }}
Environment=CHALK_TURN_URLS=turn:{{ .Domain }}:3478
Environment=CHALK_TURN_SECRET=${CHALK_TURN_SECRET}
{{- end }}

[Service]
Restart=always

[Install]
WantedBy=default.target
