# internal/chalkctl/templates/chalk-caddy.container.tmpl
# chalk -- Caddy reverse proxy + HTTP-01 TLS (generated by chalkctl init).
[Unit]
Description=chalk Caddy
Wants=network-online.target
After=network-online.target chalkd.service
Requires=chalkd.service

[Container]
ContainerName=caddy
Image=docker.io/library/caddy:{{ .CaddyTag }}
Network=chalk.network
PublishPort=80:80
PublishPort=443:443
Volume=/etc/chalk/caddy/Caddyfile:/etc/caddy/Caddyfile:ro,Z
Volume=chalk-caddy-data.volume:/data
Volume=chalk-caddy-config.volume:/config

[Service]
Restart=always

[Install]
WantedBy=default.target
